# com.ipc.bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ipc.bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VENV_PYTHON__</string>
        <string>__PROJECT_DIR__/bot.py</string>
    </array>

    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Europe/Prague</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>BOT_TOKEN</key>
        <string>__BOT_TOKEN__</string>
        <key>CHAT_ID</key>
        <string>__CHAT_ID__</string>
        <key>APP_TRV</key>
        <string>__APP_TRV__</string>
        <key>APP_WRK</key>
        <string>__APP_WRK__</string>
        <key>APP_OLD</key>
        <string>__APP_OLD__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__PROJECT_DIR__</string>

    <key>StandardOutPath</key>
    <string>__PROJECT_DIR__/ipc_monitor.log</string>

    <key>StandardErrorPath</key>
    <string>__PROJECT_DIR__/ipc_monitor.log</string>
</dict>
</plist>
